# config/supervisor.cookbook.ini
[program:cookbook]
command = /home/cookbook/config/gunicorn.sh                    ; Command to start app
user = root

stdout_logfile = /var/log/cookbook.log
redirect_stderr = true                                                ; Save std